# etc/systemd/system/boot-notify.timer
[Unit]
Description="notify people when the server is online"

[Timer]
OnBootSec=2min
AccuracySec=30sec

[Install]
WantedBy=default.target
