# etc/systemd/system/safe-poweroff.service
[Unit]
Description=poweroff

[Install]
WantedBy=default.target

[Service]
Type=simple
ExecStart=/usr/local/bin/safe-poweroff.sh
