# etc/systemd/system/snap-sync.service
[Unit]
Description=Run snap-sync backup

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/systemd-inhibit --who="snap-sync" --why="snap-sync is taking a backup" snap-sync --UUID f0a2893d-bb2f-4b41-ba46-2089fc40266a --subvolid 260 -c "root home" -q -n